# data/hfd-service.conf
description "hfd service for feedback support"

start on started dbus
stop on stopped dbus

exec hfd-service
